# deploy/claw-daily-analysis.service
[Unit]
Description=Claw 数字员工 - 每日分析一次性任务
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/AI.financial.report
Environment="PATH=/home/ubuntu/miniconda3/envs/claw/bin:/usr/bin:/bin"
Environment="DEEPSEEK_API_KEY=your_deepseek_api_key_here"
ExecStart=/home/ubuntu/miniconda3/envs/claw/bin/python run_daily.py
StandardOutput=append:/home/ubuntu/AI.financial.report/logs/daily.log
StandardError=append:/home/ubuntu/AI.financial.report/logs/daily_error.log
